# roles/keepalived/files/keepalived.target
[Unit]
Description=Collection of services that get started when keepalived is MASTER
PartOf=keepalived.service
